# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
;username = bboot
;password = bboot

[inet_http_server]
port = 0.0.0.0:5555
username = bboot
password = bboot

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=20
loglevel=error
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=ubuntu

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
prompt=web_blog

[include]
files = /etc/supervisor/conf.d/*.conf

[program:app_blog]
command=python manage.py gunicorn ; 启动命令
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user=bboot
priority=100
redirect_stderr=true
stdout_logfile=/var/log/supervisord/app_blog.log
directory=%(ENV_HOME)s/pythondir/gitpath/python_project
stdout_logfile_maxbytes=200MB

